# package/DEBIAN/prerm
#!/bin/bash

set -e

systemctl disable pi-bootstrap

